# system/skeleton/etc/profile.d/xdg.sh
export XDG_RUNTIME_DIR=${XDG_RUNTIME_DIR:-/var/run}
