# 61-usbdevice.rules
SUBSYSTEM=="android_usb",ACTION=="change",RUN+="/usr/bin/usbdevice update"
SUBSYSTEM=="udc",ACTION=="add",RUN+="/usr/bin/usbdevice update"
