# package/irqbalance/irqbalance.service
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
RuntimeDirectory=irqbalance
ExecStart=/usr/sbin/irqbalance --foreground --policyscript=/etc/irqbalance.d

[Install]
WantedBy=multi-user.target
