# board/rockchip/rk3036/fs-overlay/etc/profile.d/weston-extra.sh
# Primary screen
export WESTON_DRM_PRIMARY=HDMI-A-1

# Single screen
export WESTON_DRM_SINGLE_HEAD=1