# board/rockchip/yx/fs-overlay/usr/lib/systemd/system/rk_resize2fs.service
[Unit]
Description=Resize root partition
Before=docker.service

[Service]
ExecStart=/usr/bin/rk_resize2fs.sh

[Install]
WantedBy=multi-user.target
