# package/frecon/etc/profile.d/vt-emulation.sh
export VT_EMULATED=/var/run/vt-emulated
